# supervisor-cagecoin.conf
[supervisord]
nodaemon=true

[program:cagecoind]
command=su -c "cagecoind -conf=/home/crypto/.cagecoin/cagecoin.conf" crypto
autorestart=true
